# supervisord.sample.ini
[program:webmail_log_fetcher]
directory=/home/log_fetcher
command=/home/log_fetcher/env/bin/python /home/log_fetcher/main.py
log_stdout=true
log_stderr=true
logfile=/tmp/log/log_fetcher.log
logfile_maxbytes=10MB
logfile_backups=2